# DeepfakeBench/training/wma/wma-backend.service
# wma-backend.service
# Systemd service file for running the WMA streaming backend.

[Unit]
Description=WMA Streaming Backend Service
After=network.target

[Service]
User=roee
Group=roee

# Set the working directory to the PARENT of the 'wma' package.
WorkingDirectory=/home/roee/repos/Effort-AIGI-Detection-Fork/DeepfakeBench/training

# Standard data directory
Environment="DATA_DIRECTORY=/var/data/wma_storage"

# Execute the 'server' module inside the 'wma' package.
# This correctly handles all relative and absolute package imports.
ExecStart=/home/roee/repos/Effort-AIGI-Detection-Fork/venv/bin/python -m wma.server

Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target